# orderstream/deploy/orderstream-backend.service
# FILE: /home/overlord/.openclaw/workspace/VULCAN/deploy/orderstream-backend.service
# ─────────────────────────────────────────────────────────────────────────────
# Systemd unit for WAW Backend — direct host deployment (no Docker).
# Used when Coolify is not available or for bare-metal VPS setup.
#
# Install:
#   sudo cp deploy/orderstream-backend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable orderstream-backend
#   sudo systemctl start orderstream-backend
#
# Logs:
#   journalctl -u orderstream-backend -f
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=OrderStream (WOAW) FastAPI Backend
Documentation=https://github.com/overlord-empire/orderstream
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=overlord
Group=overlord

# Virtual environment Python — use the OVERLORD Empire standard venv
ExecStart=/opt/overlord/venv/bin/uvicorn \
    orderstream.backend.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 2 \
    --log-level info \
    --access-log

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/overlord/orderstream/.env

# Working directory must be the parent of orderstream.backend (the Python package root)
WorkingDirectory=/opt/overlord/orderstream

# Restart policy
Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/
ReadWritePaths=/opt/overlord/orderstream
ReadWritePaths=/tmp

# Security — drop capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orderstream-backend

# Health check via ExecStartPost
ExecStartPost=/bin/bash -c \
    'for i in 1 2 3 4 5; do curl -sf http://localhost:8000/health && break || sleep 2; done'

[Install]
WantedBy=multi-user.target