# orderstream/backend.service
[Unit]
Description=OrderStream API Backend
After=network.target

[Service]
User=overlord
WorkingDirectory=/home/overlord/.openclaw/workspace/VULCAN/orderstream
ExecStart=/usr/local/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
EnvironmentFile=/home/overlord/.openclaw/workspace/VULCAN/.env

[Install]
WantedBy=multi-user.target